# vps/mt5-terminal.service
[Unit]
Description=MetaTrader 5 Terminal (Wine + Xvfb)
After=xvfb.service
Requires=xvfb.service

[Service]
Type=simple
User=niri
Environment=DISPLAY=:99
Environment=WINEPREFIX=/home/niri/.wine-mt5
Environment=WINEARCH=win64
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/wine "/home/niri/.wine-mt5/drive_c/Program Files/MetaTrader 5/terminal64.exe" /portable
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
